# working/src/iterllm.service.template
[Unit]
Description=IterLLM Runtime Service
After=network-online.target

[Service]
# EnvironmentFile removed - values are baked in at install time

# --- Security & Isolation ---
Type=exec
PrivateMounts=yes
PrivateTmp=yes
ProtectSystem=full
NoNewPrivileges=yes
Delegate=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictSUIDSGID=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes

# --- Credentials ---
# Decrypts %h/.config/credstore.encrypted/my_api_key -> $CREDENTIALS_DIRECTORY/my_api_key
LoadCredentialEncrypted=my_api_key:%h/.config/credstore.encrypted/my_api_key

# --- Container Execution ---
ExecStart=/usr/bin/podman run \
    --name {{SERVICE_NAME}} \
    --replace \
    --rm \
    --cgroup-manager=systemd \
    --sdnotify=conmon \
    --network=slirp4netns \
    -v "$CREDENTIALS_DIRECTORY/my_api_key:/run/secrets/my_api_key:ro" \
    {{RUNTIME_IMAGE}}:latest

# Cleanup
ExecStop=/usr/bin/podman stop --ignore --time 10 {{SERVICE_NAME}}
ExecStopPost=/usr/bin/podman rm --force --ignore {{SERVICE_NAME}}

[Install]
WantedBy=default.target
